# automation/com.ldgr.budget-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Daily 9 AM budget-sync (macOS launchd). Replace REPLACE_WITH_HOME with your home
     directory (e.g. /Users/you), copy to ~/Library/LaunchAgents/, then:
       launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.ldgr.budget-sync.plist -->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.ldgr.budget-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>REPLACE_WITH_HOME/.claude/skills/budget-sync/run_budget_sync.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>9</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>StandardOutPath</key><string>REPLACE_WITH_HOME/.claude/skills/budget-sync/launchd.out.log</string>
  <key>StandardErrorPath</key><string>REPLACE_WITH_HOME/.claude/skills/budget-sync/launchd.err.log</string>
</dict>
</plist>
